# back-end/scripts/KTCScraperScript/crontab
* * * * * python /app/ktc-rankings-scraper.py >> /var/log/cron.log 2>&1
